# roles/minecraft_server_setup/templates/mc-router.service.j2
[Unit]
Description=Minecraft Reverse Proxy Service

[Service]
User={{ ansible_user_id }}
ExecStart=mc-router -mapping={{ mapping }}

[Install]
WantedBy=minecraft.target
