# roles/minecraft_setup/files/minecraft-servers.target
[Unit]
Description=Minecraft Servers

[Install]
WantedBy=multi-user.target
